# snet-backend.service
[Unit]
Description=S.Net Backend API Service
After=network.target postgresql.service redis.service
Wants=postgresql.service redis.service

[Service]
Type=simple
User=fathur
WorkingDirectory=/home/fathur/serverdev
Environment="NODE_ENV=production"
ExecStart=/home/fathur/.bun/bin/bun run src/index.ts
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=snet-backend

# Security settings
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
